# supervisord.conf
[supervisord]
nodaemon=true

[program:sshd]
command=/usr/sbin/sshd -D

[program:docker]
command=docker -d 

[program:fig]
command=cd /root && fige && docker2etcd && ruby ping.rb 

